# embedded-device/buildroot-external/board/stm32mp1/rootfs_overlay/etc/systemd/system/sensehat.service
[Unit]
Description=Sense Hat service

[Service]
ExecStart=
ExecStart=/bin/bash /etc/sensehat.sh

[Install]
WantedBy=default.target
